# docker-php/supervisor/etc/conf.d - 副本/micro-api-datacenter.conf
[program:micro-api-datacenter]
process_name=%(program_name)s
command=/bin/bash -c "rm -rf /var/www/html/micro-api-datacenter/runtime/container && php /var/www/html/micro-api-datacenter/bin/hyperf.php start"
autostart=true
autorestart=false
redirect_stderr=true
stdout_logfile=/var/log/supervisor/micro-api-datacenter.log